# launchd/plists/com.agrisafe.sync-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Smart orchestrator — probes all 25 sources, only runs jobs with new data
  (Daily 03:00 local — probes first, skips unchanged sources. Weekly jobs only fire on Sunday.)

  Generated by launchd/generate-plists.js from launchd/jobs.json.
  DO NOT EDIT BY HAND — re-run the generator.

  The three REPLACE_ME placeholders are filled in by launchd/install.sh
  when you install on the Mac.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agrisafe.sync-orchestrator</string>

  <key>WorkingDirectory</key>
  <string>REPLACE_ME_ABSOLUTE_PATH_TO_REPO</string>

  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_ME_ABSOLUTE_PATH_TO_NODE</string>
    <string>--env-file=.env.local</string>
    <string>node_modules/.bin/tsx</string>
    <string>src/scripts/cron/run-job.ts</string>
    <string>sync-orchestrator</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/REPLACE_ME_USERNAME/Library/Logs/AgriSafe/sync-orchestrator.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/REPLACE_ME_USERNAME/Library/Logs/AgriSafe/sync-orchestrator.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
